# prometheus-sensor-exporter.service
[Unit]
Description=Prometheus exporter for I²C temperature and humidity sensors
After=network.target

[Service]
User=prometheus
Group=i2c
EnvironmentFile=/etc/default/prometheus-sensor-exporter
ExecStart=/usr/bin/prometheus-sensor-exporter $ARGS
Restart=on-failure
LimitFSIZE=0
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
